# src/python/supervisord.conf
# boilerplate crap

[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


# our programs

[program:ipfs]
command=bash /app/scripts/start-ipfs.sh

[program:bacalhau-serve]
command=bash /app/scripts/start-bacalhau.sh
# crashloops until ipfs is ready, so restart a decent number of times
startretries=100

[program:modicum-resource-provider]
command=/usr/local/bin/modicum startRP --path /app/0_experiments/demo/ --index 1
environment=BACALHAU_API_HOST=localhost